# app/launchd/com.hxguardian.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hxguardian.server</string>

    <!-- NOTE: This file is a reference template only.
         install.sh generates the real plist inline with UserName substituted. -->
    <key>UserName</key>
    <string>_ACTUAL_USER_</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Library/Application Support/hxguardian/bin/run-hxg-server</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Library/Logs/hxguardian-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Logs/hxguardian-server-error.log</string>
</dict>
</plist>
